# machines/eb-mailslurper/etc/systemd/system/mailslurper.service
[Unit]
Description=mailslurper
After=network-online.target

[Service]
User=mailslurper
Group=mailslurper
WorkingDirectory=/home/mailslurper/app/
ExecStart=/home/mailslurper/app/mailslurper
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
